# conf/supervisor.conf
[supervisord]
nodaemon=true
user=root

[inet_http_server]
port=0.0.0.0:8000
#username=qtluser # Basic auth username
#password=qtlpass # Basic auth password

[program:rserve]
command=/usr/local/bin/Rscript /app/qtl2rest/run.R
priority=1
autostart=true
autorestart=false
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
